# config/systemd/alexandria-sync.service
[Unit]
Description=Alexandria Sneakernet Import Verification Service
Documentation=file:///opt/alexandria/docs/maintenance-guide.md
After=local-fs.target

[Service]
Type=oneshot
User=alexandria
Group=alexandria
WorkingDirectory=/opt/alexandria

# Verify the current ZIM archive integrity after a sneakernet import
ExecStart=/opt/alexandria/scripts/verify_zim.sh /mnt/pve/alexandria/zim

# Harden the service unit itself
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/mnt/pve/alexandria/zim
PrivateTmp=true
PrivateDevices=true
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
SystemCallFilter=@system-service
CapabilityBoundingSet=

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=alexandria-sync

[Install]
WantedBy=multi-user.target
